# dtrack/dtrack_gateway.service
[Unit]
Description=DTrack Gateway
After=multi-user.target

[Service]
Type=simple
Restart=always
RestartSec=30
WorkingDirectory=/home/dtrack/dtrack_gateway
ExecStart=/usr/bin/python3 /home/dtrack/dtrack_gateway/gateway.py

[Install]
WantedBy=network-online.target
